# netatalk/supervisord.conf
[supervisord]
nodaemon=true

[program:timemachine]
command=/start_netatalk.sh

[program:avahi]
command=avahi-daemon -D
